# notes/agua/private/bin/test/node/mysql/node/init/rabbitjs.conf
description "Run rabbit.js"
author "Stuart Young"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

script

    sleep 10;
    /usr/bin/nodejs INSTALLDIR/html/exchange/rabbitjs/example/socketio/server.js 2>&1 >> INSTALLDIR/log/rabbitjs.log

end script
